# supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisord.log
pidfile=/var/run/supervisord.pid

[program:netrekd]
command=/opt/netrek/bin/netrekd
directory=/opt/netrek
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
priority=10

[program:ws-proxy]
command=node /opt/ws-proxy/index.js
autostart=true
autorestart=true
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0
environment=NETREK_HOST="127.0.0.1",NETREK_PORT="%(ENV_NETREK_PORT)s",WS_PORT="%(ENV_WS_PORT)s",STATIC_DIR="/opt/web-client",PORTAL_DIR="/opt/portal"
priority=20
